# profile.d/bash_aliases.sh
#!/bin/bash

# shorten prompt
alias sprompt="PROMPT_DIRTRIM=3"
